# docker/dockerfile/supervisor.d/mixfilecli.ini
[program:mixfilecli]
directory = /alist-mixfileweb/mixfilecli/
command = java -jar mixfile-cli.jar
autostart = true
stopsignal = INT
startsecs = 20
autorestart = true
startretries = 3
redirect_stderr=true
stdout_logfile_maxbytes = 10MB
stdout_logfile_backups = 10
stdout_logfile = /alist-mixfileweb/mixfilecli/mixfilecli.log