# docker/dockerfile/supervisor.d/mixfile_webdav.ini
[program:mixfile_webdav]
directory = /alist-mixfileweb/mixfile_webdav/
command = node app.js
autostart = true
stopsignal = INT
startsecs = 20
autorestart = true
startretries = 3
redirect_stderr=true
stdout_logfile_maxbytes = 10MB
stdout_logfile_backups = 10
stdout_logfile = /alist-mixfileweb/mixfile_webdav/mixfile_webdav.log